# cfg/rpi/etc/systemd/system/getty@tty1.service.d/override.conf
[Service]
ExecStart=
ExecStart=-/usr/bin/startx -- -nocursor
User=pi
StandardInput=tty
StandardOutput=tty
